# launchd/com.nielsmadan.hidutil-capslock-to-f18.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nielsmadan.hidutil-capslock-to-f18</string>

    <!--
      Remap Caps Lock (HID 0x700000039) to F18 (HID 0x70000006D) at the
      HID-driver level. F18 is then picked up by Hammerspoon as the
      modal-entry trigger for the window-management leader sequence.

      hidutil's UserKeyMapping is session-scoped (lost across reboot
      and across full logouts), so this LaunchAgent re-applies it at
      every login via RunAtLoad. Symlinked into ~/Library/LaunchAgents/
      by install.sh.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/hidutil</string>
        <string>property</string>
        <string>--set</string>
        <string>{"UserKeyMapping":[{"HIDKeyboardModifierMappingSrc":0x700000039,"HIDKeyboardModifierMappingDst":0x70000006D}]}</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
